# backend-laravel/deploy/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
;  PANDA — Supervisor process definitions (production)
;
;  Drop this in /etc/supervisor/conf.d/panda.conf and then:
;
;     sudo supervisorctl reread
;     sudo supervisorctl update
;     sudo supervisorctl status
;
;  Three processes:
;    panda-queue     queue workers (broadcast jobs, mail, notifications)
;    panda-reverb    Reverb websocket server (Pusher protocol on :8080)
;    panda-scheduler artisan schedule:work — keeps cron jobs running
;
;  Replace /var/www/panda with your actual deploy path.
; ─────────────────────────────────────────────────────────────────────────────

[program:panda-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/panda/artisan queue:work redis --queue=default,broadcasting --sleep=1 --tries=3 --max-time=3600 --backoff=5
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/panda-queue.log
stopwaitsecs=3600
stopasgroup=true
killasgroup=true

[program:panda-reverb]
process_name=%(program_name)s
command=php /var/www/panda/artisan reverb:start --host=0.0.0.0 --port=8080 --no-interaction
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/panda-reverb.log
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:panda-scheduler]
process_name=%(program_name)s
command=php /var/www/panda/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/panda-scheduler.log
